# remote_printer.timer
[Unit]
Description=Remote Printer Timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target
